# rootfs/overlay/etc/inittab
::sysinit:/etc/init.d/rcS
::respawn:/sbin/getty -L console 0 vt100
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
